# .config/autostart/blueman-applet.desktop
[Desktop Entry]
Name=Blueman Applet
GenericName=Bluetooth Menu
Exec=blueman-applet
Type=Application
StartupNotify=false
