# .config/autostart/nm-applet.desktop
[Desktop Entry]
Name=NM Applet
GenericName=Wifi Menu
Exec=nm-applet
Type=Application
StartupNotify=false
